# scripts/gunicorn.conf
[program:fludetector]
command = /home/deploy/fludetector/scripts/serve.sh
directory = /home/deploy/fludetector
user = deploy
